# battery-notify/.config/systemd/user/battery-notify.service
[Unit]
Description=Notification of low battery

[Service]
Type=simple
ExecStart=%h/bin/battery-notify

[Install]
WantedBy=default.target
